# dist/fanctl.service
# /etc/systemd/system/fanctl.service

[Unit]
Description=fan control service

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/fanctl -loglevel 1
Restart=always
RestartSec=3000ms

[Install]
WantedBy=multi-user.target

# https://www.freedesktop.org/software/systemd/man/systemd.service.html
